# build/com.local.quoteDownloader.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>com.local.quoteDownloader</string>

        <key>Program</key>
        <string>/Users/wolt/Documents/01 Workplace William/Quote Downloader/ftpDownloader/build/quote-downloader</string>

        <key>CPU</key>
        <string>15</string>

        <key>RunAtLoad</key>
        <true/>

        <key>KeepAlive</key>
        <true/>
	</dict>
</plist>